# src/crontab
SHELL=/bin/bash

* * * * * root      /home/ctf/backup.sh
